# capture_only/online/etc/systemd/system/atn_capture.service
[Unit]
Description=Capture ATN traffic
Wants=network-online.target
After=network-online.target

[Service]
# File with environment variables used in the Bash scripts
EnvironmentFile=/etc/sysconfig/atn_capture.env

# Run cleanup before starting; tolerate failures so service start is not blocked
ExecStartPre=-/opt/muac/atn_capture_cleanup.sh
ExecStart=/opt/muac/atn_capture.sh
Type=simple
TimeoutStopSec=5s
# Lifecycle managed by Pacemaker
Restart=no

[Install]
WantedBy=multi-user.target
